# RasPi/etc/supervisor/conf.d/shoppi.conf
[program:ShopPi]
command=/home/pi/ShopPi/ShopPi
user=pi
directory=/home/pi/ShopPi
