# RasPi/etc/supervisor/conf.d/shoppi.conf
[program:ShopPi]
command=/home/runner/ShopPi/ShopPi
user=runner
directory=/home/runner/ShopPi
