# ansible-role/templates/etc/supervisor/conf.d/shoppi.conf.j2
[program:ShopPi]
command={{ runner_home_shoppi }}/ShopPi
user={{ runner_user }}
directory={{ runner_home_shoppi }}
